# rdzbridge/root/etc/uci-defaults/98-rdzbridge-finalize
#!/bin/sh
set -eu

WRAPPER_LOG='/etc/rdzbridge-wrapper.log'
READY='/etc/rdzbridge-ready.json'
SELF='/etc/uci-defaults/98-rdzbridge-finalize'

exec >>"$WRAPPER_LOG" 2>&1

echo "RDZ wrapper start: $(date -Iseconds 2>/dev/null || date)"

[ -x /usr/libexec/rdzbridge-finalize ] || {
	echo 'ERROR: permanent RDZ finalizer is missing or not executable'
	exit 1
}

# Enable the recovery path before applying the main configuration. If the
# finalizer fails, this wrapper remains present and recovery retries next boot.
/etc/init.d/rdzbridge-recovery enable >/dev/null 2>&1 || true

/usr/libexec/rdzbridge-finalize

[ -s "$READY" ] || {
	echo 'ERROR: finalizer returned without creating the ready marker'
	exit 1
}

STATUS="$(jsonfilter -i "$READY" -e '@.status' 2>/dev/null || true)"
[ "$STATUS" = ready ] || {
	echo "ERROR: unexpected ready-marker status: ${STATUS:-missing}"
	cat "$READY" 2>/dev/null || true
	exit 1
}

rm -f "$SELF"
echo "RDZ wrapper complete: $(cat "$READY")"
sync
exit 0
