# rdzbridge/root/etc/uci-defaults/98-rdzbridge-finalize
#!/bin/sh

# Enable a permanent recovery path before applying the configuration. If the
# finalizer fails for any transient reason, the recovery service retries on the
# next boot instead of leaving the appliance inaccessible.
/etc/init.d/rdzbridge-recovery enable >/dev/null 2>&1 || true
exec /usr/libexec/rdzbridge-finalize
