# rdzbridge/root/etc/uci-defaults/98-rdzbridge-finalize
#!/bin/sh

# Remove this transient uci-default before starting the permanent finalizer.
# This makes repeated/manual execution safe and prevents CI or recovery logic
# from observing a stale first-boot script while configuration is still active.
rm -f "$0" 2>/dev/null || true

# Enable a permanent recovery path before applying the configuration.